# home/.config/autostart/Disable Capslock.desktop
[Desktop Entry]
Type=Application
Name=Disable Capslock
Description=
Exec=setxkbmap -option ctrl:nocaps
Icon=keyboard